# base_external/rootfs_overlay/etc/init.d/S98lddmodules
#!/bin/sh
### BEGIN INIT INFO
# Provides:          lddmodules
# Required-Start:    $all
# Required-Stop:      0 6
# Default-Start:     S
# Default-Stop:      0 6
# Short-Description: Load LDD modules at boot
### END INIT INFO

case "$1" in
  start)
    echo "Loading kernel modules..."
    modprobe scull || { echo "Failed to load scull"; exit 1; }
    modprobe faulty || { echo "Failed to load faulty"; rmmod scull; exit 1; }
    modprobe hello || { echo "Failed to load hello"; rmmod faulty; rmmod scull; exit 1; }
    echo "All modules loaded"
    ;;
  stop)
    echo "Unloading kernel modules..."
    rmmod hello || echo "Failed to unload hello"
    rmmod faulty || echo "Failed to unload faulty"
    rmmod scull || echo "Failed to unload scull"
    echo "Module unload complete"
    ;;
  *)
    echo "Usage: $0 {start|stop}"
    exit 1
    ;;
esac
